# modemrebooter.timer
[Unit]
Description=modemrebooter timer

[Timer]
OnUnitActiveSec=120s
OnBootSec=3min

[Install]
WantedBy=timers.target